# docker/s6/web/run
#!/command/with-contenv sh

caddy run --config /etc/caddy/Caddyfile --adapter caddyfile
