# docker/s6/api/run
#!/command/with-contenv sh

cd /app/api
./mishosapi
